# debian/cron.daily
#!/bin/bash
export PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin:$PATH
. /etc/maldetect/conf.maldet

if [ -r "/etc/default/maldetect" ]; then
    source /etc/default/maldetect
fi

if [ ! "${STDDOCROOT}" ]; then
    STDDOCROOT="/var/www/html /usr/local/apache/htdocs /var/www/ /srv/www/ /home/?/public_html/ /home/?/domains/?/public_html/"
fi

if [ ! "${CRON_DAYS}" ]; then
    CRON_DAYS="2"
fi

find=`which find 2> /dev/null`
if [ "$find" ]; then
    # prune any quarantine/session/tmp data older than 14 days
    tmpdirs="/tmp/maldetect /var/lib//maldetect/sess /var/lib/maldetect/quarantine /var/lib/maldetect/pub"
    for dir in $tmpdirs; do
     if [ -d "$dir" ]; then
      $find $dir -type f -mtime +14 -print0 | xargs -0 rm -f >> /dev/null 2>&1
     fi
    done
fi

# if were running inotify monitoring, send daily hit summary
if [ "$(ps -A --user root -o "comm" | grep inotifywait)" ]; then
        maldet --alert-daily >> /dev/null 2>&1
else
    # scan the last n days of file changes
    for DIR in ${STDDOCROOT}; do
        maldet -b -r ${DIR} ${CRON_DAYS} >> /dev/null 2>&1
    done
fi
